# src/minecraftd-save@.timer
[Unit]
Description=Daily save of Minecraft server %i

[Timer]
OnCalendar=daily
Persistent=true


[Install]
WantedBy=multi-user.target

