# src/minecraftd-backup@.service
[Unit]
Description=Backup %I Minceraft server

[Service]
Type=oneshot
EnvironmentFile=-/etc/minecraft/%i
ExecStart=/usr/bin/minecraftctl backup %i
